# scripts/native/systemd/zip-serial-bridge.service
[Unit]
Description=ZIP Robot Serial Bridge (ROS 2 Arduino Communication)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i

# Environment setup
Environment="ROS_DOMAIN_ID=0"
Environment="ROS_DISTRO=humble"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/ros/humble/bin:/home/%i/.local/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/%i/.local/lib/python3.10/site-packages:/usr/lib/python3.10/site-packages"

# Source ROS 2 and workspace
ExecStartPre=/bin/bash -c 'source /opt/ros/humble/setup.bash; source /home/%i/ros2_ws/install/setup.bash'

# Launch serial bridge node
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash; source /home/%i/ros2_ws/install/setup.bash; ros2 launch zip_control serial_bridge.launch.py'

# Restart policy
Restart=on-failure
RestartSec=10

# Resource limits
MemoryLimit=256M
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zip-serial-bridge

# Device access (for serial ports)
SupplementaryGroups=dialout tty

[Install]
WantedBy=multi-user.target
